# deploy/systemd/regionlens-backup.service
# Резервное копирование PostgreSQL RegionLens — одна прогонка скрипта pg_backup.sh.
#
# Запускается по таймеру regionlens-backup.timer (ежедневно) либо вручную:
#   systemctl start regionlens-backup.service
# Скрипт делает pg_dump из контейнера postgres в gzip с ротацией (KEEP_BACKUPS).
#
# Пути ниже — примерная раскладка развёртывания (/opt/regionlens, compose-стек рядом,
# переменные окружения в /etc/regionlens.env) — поправьте под свой сервер.
#
# Установка:
#   sudo cp deploy/systemd/regionlens-backup.* /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now regionlens-backup.timer

[Unit]
Description=RegionLens: резервная копия PostgreSQL (pg_dump с ротацией)
# Бэкап осмыслен только при поднятом compose-стеке (контейнер postgres).
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
User=regionlens
Group=regionlens
WorkingDirectory=/opt/regionlens
EnvironmentFile=/etc/regionlens.env
# Куда и сколько хранить — переопределяемо (по умолчанию data/backups/postgres, 7 копий).
Environment=BACKUP_DIR=/opt/regionlens/data/backups/postgres
Environment=KEEP_BACKUPS=7
ExecStart=/bin/bash /opt/regionlens/deploy/backup/pg_backup.sh
TimeoutStartSec=15min
# Журнал — в systemd journal: journalctl -u regionlens-backup.service
StandardOutput=journal
StandardError=journal
